# hubot-worker.conf
[program:hubot]
command = /path/to/bin/hubot -a irc --name MyName

autostart = true
autorestart = true
user = username
directory = /path/to/hubot
redirect_stderr = true
stdout_logfile = /dev/null
